# com.contentos.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.contentos.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/rajvi/Desktop/ContentOS/run-tunnel.sh</string>
  </array>

  <!-- Start on login, and relaunch automatically if it ever exits/crashes -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/rajvi/Desktop/ContentOS/tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rajvi/Desktop/ContentOS/tunnel.log</string>
</dict>
</plist>
